# ops/systemd-user/cryptoalpha-duty-check.service
[Unit]
Description=CryptoAlpha Duty Check
After=network-online.target cryptoalpha-backend.service
Wants=network-online.target cryptoalpha-backend.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -lc 'set -euo pipefail; curl -fsS http://127.0.0.1:8000/api/ops/duty-check'
