# ops/systemd-user/cryptoalpha-recommender-events.timer
[Unit]
Description=CryptoAlpha RL Recommender Events Digest Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=cryptoalpha-recommender-events.service

[Install]
WantedBy=timers.target
