# ops/systemd-user/cryptoalpha-retention.timer
[Unit]
Description=CryptoAlpha Retention Timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
RandomizedDelaySec=2min
Persistent=true
Unit=cryptoalpha-retention.service

[Install]
WantedBy=timers.target
